# s6-rc.d/php-fpm/run
#!/command/execlineb -P
s6-setuidgid http
php-fpm8.2 --nodaemonize
